# bootc-custom-image/bootc-fetch-apply-updates.timer.d/updates.conf
[Timer]
OnBootSec= 
OnBootSec=5m 
OnUnitInactiveSec=5m
OnUnitActiveSec=5min
OnActiveSec=5m
OnCalendar=*:0/5
Persistent=true
